# etc/supervisord.d/hubot.ini
[program:hubot]
user=hubot
directory=/home/hubot/
command=/home/hubot/bin/hubot -n $BOT_NAME -a rocketchat
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = true
